# autostart/xfce4-screensaver.desktop
[Desktop Entry]
Exec=/usr/bin/xfce4-screensaver
Name=xfce4-screensaver
Comment=The Xfce 4 Screensaver
Type=Application
