# autostart/xfce4-power-manager.desktop
[Desktop Entry]
Exec=/usr/bin/xfce4-power-manager
Name=xfce4-power-manager
Comment=The Xfce 4 Power manager
Type=Application
